# scadalts/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/sbin/mysqld --datadir=/var/lib/mysql --user=mysql
autostart=true
autorestart=true
priority=5

[program:init]
command=/init.sh
autostart=true
autorestart=false
priority=7
startsecs=0
startretries=1
stdout_logfile=/var/log/supervisor/init.log
stderr_logfile=/var/log/supervisor/init.err

[program:tomcat]
command=/bin/bash -c "sleep 10 && /usr/local/tomcat/bin/catalina.sh run"
autostart=true
autorestart=true
priority=10

[program:wazuh-agent]
command=/bin/bash -c "if getent hosts wazuh >/dev/null 2>&1; then /var/ossec/bin/wazuh-control start; else echo '[scadalts] Wazuh not in DNS, skipping agent start'; fi; sleep infinity"
autostart=true
autorestart=false
priority=20
stdout_logfile=/var/log/supervisor/wazuh-agent.log
stderr_logfile=/var/log/supervisor/wazuh-agent.err
